# mkosi.extra/usr/lib/systemd/system/systemd-sysupdated.service.d/kde-linux-tar-options.conf
# SPDX-License-Identifier: GPL-2.0-only OR GPL-3.0-only OR LicenseRef-KDE-Accepted-GPL
# SPDX-FileCopyrightText: 2024 Lasath Fernando <devel@lasath.org>

[Service]

# Our images are compressed with zst, but systemd-pull does not support
# zstd compression yet:
# https://github.com/systemd/systemd/issues/33905
#
# However, it turns out GNU tar itself reads environment variables:
# https://www.gnu.org/software/tar/manual/tar.html#TAR_005fOPTIONS
Environment="TAR_OPTIONS=--zstd"
